# src/main/assembly/dist/install/dku-todo-db.service
[Unit]
Description=${project.name} Service

[Service]
ExecStart=/opt/dans.knaw.nl/dku-todo-db/bin/dku-todo-db server /etc/opt/dans.knaw.nl/dku-todo-db/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143

User=dku-todo-db
Group=dku-todo-db

[Install]
WantedBy=multi-user.target
